# rootfs/etc/s6-overlay/s6-rc.d/goofys-wait/up
#!/usr/bin/env sh

while ! df | grep -qE '^[[:alnum:]]{7}.*/data$'; do
    sleep 1
done
