# rootfs/etc/s6-overlay/s6-rc.d/mountpoint-wait/up
#!/usr/bin/env sh

while ! df --output=target,source | grep -qE '^/data\s+mountpoint-s3$'; do
    sleep 1
done
