# rootfs/etc/s6-overlay/s6-rc.d/nbd-init/up
#!/usr/bin/env sh

mkdir -p /mnt/zerofs-9p/.nbd/
truncate -s 1G /mnt/zerofs-9p/.nbd/storage
