# rootfs/etc/s6-overlay/s6-rc.d/nbd-mount/up
#!/usr/bin/env sh

mkdir -p /mnt/zerofs-nbd4

mount /dev/nbd4 /mnt/zerofs-nbd4
