# rootfs/etc/s6-overlay/s6-rc.d/s3fs-wait/up
#!/usr/bin/env sh

while ! df --output=target,source | grep -qE '^/data\s+s3fs$'; do
    sleep 1
done
