# rootfs/etc/s6-overlay/s6-rc.d/zerofs-run/run
#!/usr/bin/env sh

exec zerofs run --config="/etc/zerofs/config.toml"
