# rootfs/etc/s6-overlay/s6-rc.d/juicefs-mount/run
#!/usr/bin/env sh

mkdir -p \
    /data/

exec juicefs mount \
    --no-usage-report \
    sqlite3://myjfs.db \
    /data/ \
    -d
