# rootfs/etc/s6-overlay/s6-rc.d/nfs-configure/up
#!/usr/bin/env sh

exportfs -r

rpc.nfsd --no-nfs-version 3 --nfs-version 4.0 --nfs-version 4.1 --nfs-version 4.2
